# AuroraQ/scripts/systemd/auroaq-trading.service
[Unit]
Description=AuroraQ Trading System
Documentation=https://github.com/your-org/AuroraQ
After=network.target auroaq-sentiment.service
Wants=network.target
Requires=auroaq-sentiment.service

[Service]
Type=simple
User=auroaq
Group=auroaq
WorkingDirectory=/opt/AuroraQ

# Environment variables
Environment=PYTHONPATH=/opt/AuroraQ
Environment=AURORA_ENV=production
Environment=AURORA_MODE=live
EnvironmentFile=-/opt/AuroraQ/.env

# Start command
ExecStart=/opt/AuroraQ/venv/bin/python -m SharedCore.main_trading_system
ExecReload=/bin/kill -HUP $MAINPID

# Pre-start checks
ExecStartPre=/bin/bash -c 'curl -f http://localhost:8000/health || exit 1'
ExecStartPre=/opt/AuroraQ/venv/bin/python /opt/AuroraQ/scripts/production_readiness_check.py --quick

# Restart policy
Restart=on-failure
RestartSec=30
StartLimitInterval=600
StartLimitBurst=3

# Resource limits
MemoryLimit=4G
CPUQuota=300%

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/AuroraQ/logs /opt/AuroraQ/SharedCore/data_storage /opt/AuroraQ/backups /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auroaq-trading

# Graceful shutdown
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target