# AuroraQ/scripts/systemd/auroaq-monitor.service
[Unit]
Description=AuroraQ Health Monitor
Documentation=https://github.com/your-org/AuroraQ
After=network.target auroaq-sentiment.service auroaq-trading.service
Wants=network.target

[Service]
Type=simple
User=auroaq
Group=auroaq
WorkingDirectory=/opt/AuroraQ

# Environment variables
Environment=PYTHONPATH=/opt/AuroraQ
Environment=AURORA_ENV=production

# Monitor command
ExecStart=/opt/AuroraQ/venv/bin/python /opt/AuroraQ/sentiment-service/scripts/health_monitor.py --monitor

# Restart policy
Restart=on-failure
RestartSec=60
StartLimitInterval=600
StartLimitBurst=5

# Resource limits
MemoryLimit=512M
CPUQuota=50%

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/AuroraQ/logs /tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auroaq-monitor

[Install]
WantedBy=multi-user.target