# AuroraQ/scripts/systemd/auroaq-backup.timer
[Unit]
Description=AuroraQ Backup Timer
Requires=auroaq-backup.service

[Timer]
# Run daily at 2 AM
OnCalendar=*-*-* 02:00:00
# Run if missed (e.g., system was down)
Persistent=true
# Random delay up to 10 minutes to avoid system load spikes
RandomizedDelaySec=600

[Install]
WantedBy=timers.target